# crontab
# m h  dom mon dow   command
0 5  * * *  /home/ubuntu/goalification/cron.sh
